# modules/orchestration/templates/systemd.service.j2
[Unit]
Description=OpenClaw Gateway (profile: {{ profile }}, claw-swarm managed)
After=network-online.target
Wants=network-online.target

[Service]
# default profile 使用 ~/.openclaw/，其他 profile 使用 ~/.openclaw-<profile>/
# 使用 `-` 前缀表示文件不存在时忽略
EnvironmentFile=-{{ claw_user_home }}/.openclaw/runtime/env
EnvironmentFile=-{{ claw_user_home }}/.openclaw-{{ profile }}/runtime/env
ExecStart={{ node_bin }} {{ openclaw_bin }} --profile {{ profile }} gateway run --port {{ port }} --bind loopback
Restart=always
RestartSec=5
TimeoutStopSec=30
TimeoutStartSec=30
SuccessExitStatus=0 143
KillMode=control-group

[Install]
WantedBy=default.target
